# dat/rc.local
#!/bin/sh -e

echo EJECUTANDO RC.LOCAL

exit 0